# etc/supervisor/production.conf
[supervisord]
pidfile=%(here)s/../../tmp/supervisord_production.pid
logfile=%(here)s/../../tmp/supervisord_production.log

[program:django-server]
command=%(here)s/../../dev/runwsgi.py production
